# etc/cron.daily/121-slaptgetupdate
#!/bin/bash

### User Variables START ###
customtoolsrootloc="/SETAPATH_NO_TRAILING_SLASH"
windowtomoveto="4" #use count starting at 1

### User Variables END ###


## Script starts here ##

set -euo pipefail
. /etc/crontimerloop.sh
. /etc/windowlabelandmove.sh

pkgdir=$(grep '^WORKINGDIR=' /etc/slapt-get/slapt-getrc | sed 's/^WORKINGDIR=//')
echo "Downloading any available system updates"
echo
if [ -f "/tmp/slapt-get-done" ]; then
  rm /tmp/slapt-get-done
fi
if ! cd "${pkgdir}"; then
  echo "unable to cd to ${pkgdir}"
  echo "aborting!"
  echo
  echo
  exit 1
fi

exclusions_modified=0
original_excludes=""
wintomv=""
sources_modified=0
original_sources=""
disabled_sources_list=""
download_error="(HTTP response code said error|\
Failure when receiving data from the peer)"

fn_main_loop () {
  fn_reset_downloads
  fn_update
  if echo "$update_output" | grep -q "Sources failed to download, correct sources and rerun --update"; then
    fn_update_error
  fi

  fn_upgrade
  if echo "$log_output" | grep -q "MD5 checksum mismatch"; then
    fn_md5_mismatch
  elif echo "$log_output" | grep -qE "Failed to download: ${download_error}" ; then
    fn_http_error
    fn_main_loop
  elif echo "$log_output" | grep -q "OTHER ERROR"; then
    echo "Not Handling OTHER ERROR"
  fi
}

fn_reset_downloads () {
  mv UpdateAll /tmp/ 2>/dev/null
  mv OrphanDeps /tmp/ 2>/dev/null
  mv 32multi /tmp/ 2>/dev/null
  slapt-get --clean
  if [ "$(ls -A | grep -v '^\.')" ]
  then
	  rm -r * 2>/dev/null
  fi
  if [ "$(ls -A | grep '^\.')" ] ; then
  	rm -r .* 2>/dev/null
  fi
  mv /tmp/UpdateAll ${pkgdir}/ 2>/dev/null
  mv /tmp/OrphanDeps ${pkgdir}/ 2>/dev/null
  mv /tmp/32multi ${pkgdir}/ 2>/dev/null

  return 0
}

fn_update () {
  update_output=$(ctl-fn_timeoutloop "slapt-get --update" 600 3 \
    "slapt-get --update" 2>&1 || true )
  rc=$?
echo "${update_output}" | tail -n 1
  if [ $rc -eq 124 ]; then
    echo "$update_output" | tail -n 1
    fn_restore_exclusions
    exit 1
  fi

  return 0
}


fn_upgrade() {
  log_output=$(ctl-fn_timeoutloop "slapt-get --upgrade" 1800 3 \
    "slapt-get --upgrade --ignore-dep --download-only" 2>&1 || true )
  rc=$?

  if [ $rc -eq 124 ]; then
    echo "$log_output" | tail -n 1
    fn_restore_exclusions
    fn_reset_downloads
    exit 1
  fi

  return 0
}


fn_disable_sources () {
  local source_list="$1"

  # Set the flag that sources have been modified
  if [ $sources_modified -eq 0 ]; then
    sources_modified=1
  fi

  # Disable each source in the list and track what was disabled
  while IFS= read -r source; do
    if [ -n "$source" ]; then
      echo "Temporarily disabling source: $source"
      source_with_slash="${source%/}/"

      # Check if this source exists and disable it
      if grep -q "^SOURCE=${source_with_slash}$" /etc/slapt-get/slapt-getrc; then
        # Add to our list of disabled sources for later restoration
        if [ -z "$disabled_sources_list" ]; then
          disabled_sources_list="$source_with_slash"
        else
          disabled_sources_list="${disabled_sources_list} $source_with_slash"
        fi

        # Disable the source
        sed -i "s|^SOURCE=${source_with_slash}|#DISABLED=${source_with_slash}|g" /etc/slapt-get/slapt-getrc
      else
        echo "Warning: Could not find exact source match for: $source_with_slash"
      fi
    fi
  done <<< "$source_list"

  return
}


fn_update_error () {
  echo "Update sources failed, analyzing problematic sources..."

  # Extract the full URLs from lines ending with "Failed to download"
  failed_sources=$(echo "$update_output" | grep "Failed to download" | grep -oP '\[https?://[^\]]+' | sed 's/\[//')

  if [ -n "$failed_sources" ]; then
    echo "Found problematic sources:"
    echo "$failed_sources"

    fn_disable_sources "$failed_sources"

    echo "Retrying update with problematic sources disabled..."
    fn_main_loop
  else
    echo "Could not identify specific problematic sources"
    echo "Full update output:"
    echo "$update_output"
    exit 1
  fi
  return
}

fn_md5_mismatch () {
  bad_source=$(echo "$log_output" | awk '/MD5 checksum mismatch/ {print x}; {x=$0}' | grep -oP 'Get \K[^ ]+' | sed 's|.*//||')
  echo "Problematic source detected: $bad_source"

  # Convert to full URL format for consistency
  full_source="https://$bad_source"

  fn_disable_sources "$full_source"

  fn_reset_downloads
  fn_update
  fn_upgrade

  if [[ $result -eq 0 ]]; then
    echo "Upgrade successful after excluding problematic source."
  else
    echo "Upgrade failed even after excluding problematic source."
    fn_reset_downloads
    exit 1
  fi
  return
}

fn_http_error () {
  # Extract package name from error message
  error_line=$(echo "$log_output" | grep -E "Failed to download: ${download_error}")

  # Look for the pattern: "Get [URL] [PackageName] [Version]"
  if [[ "$error_line" =~ Get[[:space:]]+[^[:space:]]+[[:space:]]+([^[:space:]]+)[[:space:]] ]]; then
    local problematic_pkg="${BASH_REMATCH[1]}"

    if [ -n "$problematic_pkg" ]; then
      echo "HTTP download error detected for package: $problematic_pkg"

      # Save original exclude list if this is the first exclusion
      if [ $exclusions_modified -eq 0 ]; then
        original_excludes=$(grep "^EXCLUDE=" /etc/slapt-get/slapt-getrc | sed 's/^EXCLUDE=//')
        exclusions_modified=1
      fi

      # Get current exclusion list
      current_exclude=$(grep "^EXCLUDE=" /etc/slapt-get/slapt-getrc | sed 's/^EXCLUDE=//')

      # Add the new package to the exclude list
      new_exclude="${current_exclude},${problematic_pkg}"
      echo "Temporarily adding $problematic_pkg to exclusion list"
      sed -i "s/^EXCLUDE=.*/EXCLUDE=${new_exclude}/" /etc/slapt-get/slapt-getrc

      echo "Will retry upgrade after excluding problematic package: $problematic_pkg"
    else
      echo "Failed to identify problematic package from HTTP error"
      fn_reset_downloads
      exit 1
    fi
  else
    echo "Failed to parse HTTP error message format"
    fn_reset_downloads
    exit 1
  fi
  return
}

fn_restore_exclusions () {
  if [ $exclusions_modified -eq 1 ]; then
    echo "Restoring original exclusion list"
    sed -i "s/^EXCLUDE=.*/EXCLUDE=${original_excludes}/" /etc/slapt-get/slapt-getrc
    exclusions_modified=0
  fi

  if [ $sources_modified -eq 1 ]; then
    echo "Restoring disabled sources"
    for source in $disabled_sources_list; do
      if [ -n "$source" ]; then
        echo "Re-enabling source: $source"
        sed -i "s|^#DISABLED=${source}|SOURCE=${source}|g" /etc/slapt-get/slapt-getrc
      fi
    done
    sources_modified=0
    disabled_sources_list=""
  fi
  return
}

fn_pkgtools-update_check () {
 if slapt-get --search pkgtools  | grep -w "pkgtools" | grep -v slapt-get | grep inst=no ; then
    upgradepkgloc="/sbin/upgradepkg"
	  echo "New pkgtools available,"
  	echo "updating to new pkgtools"
	  slapt-get --install pkgtools -y > /dev/null 2>&1
    echo "backing up new ${upgradepkgloc}"
    cp -a ${upgradepkgloc} ${customtoolsrootloc}${upgradepkgloc}-new
    echo "replacing new ${upgradepkgloc} with customised version"
    cp -a "${customtoolsrootloc}${upgradepkgloc}" "${upgradepkgloc}"
    echo "removing pkgtools from ${pkgdir}"
    find "${pkgdir}" -type f -name 'pkgtools*' -delete
  fi
}

fn_java_check () {
  if downloaded_java=$(ls -d openjdk* 2>/dev/null | grep -o '[0-9]\+'); then
    installed_java=$(java -version 2>&1 | awk -F[\".] '/version/ {print $2}')
    if [ "${installed_java}" -ge "${downloaded_java}" ] ; then rm -rf ./openjdk* ; fi
  fi
}

fn_mozjs-update_check () {
  # Get all mozjs packages and their installation status
  mozjs_output=$(slapt-get --search mozjs)

  # Extract version numbers from installed and available packages
  installed_mozjs_version=""
  available_mozjs_versions=()

  while IFS= read -r line; do
	  if [[ $line =~ mozjs([0-9]+)-[0-9]+\.[0-9]+\.[0-9]+.*\[inst=yes\] ]]; then
  		installed_mozjs_version="${BASH_REMATCH[1]}"
	  elif [[ $line =~ mozjs([0-9]+)-[0-9]+\.[0-9]+\.[0-9]+.*\[inst=no\] ]]; then
  		available_mozjs_versions+=("${BASH_REMATCH[1]}")
	  fi
  done <<< "$mozjs_output"

  # Check if there's a newer version available
  newer_mozjs_available=false
  if [[ -n "$installed_mozjs_version" ]]; then
  	for mozjs_version in "${available_mozjs_versions[@]}"; do
  		if (( mozjs_version > installed_mozjs_version )); then
  			newer_mozjs_available=true
  			break
  		fi
  	done
  fi

  if $newer_mozjs_available; then
  	echo "New MozJS available,"
  	echo "recommend uninstalling the old and installing the new one"
  	echo
  	echo
  fi
}

fn_main_loop
fn_restore_exclusions
fn_pkgtools-update_check
fn_java_check
fn_mozjs-update_check

echo
excludes="$(ls | grep -v 'package_data' | grep -v 'UpdateAll' | grep -v 'OrphanDeps' | grep -v '32multi')"
if [[ ${excludes} ]] ; then
  if [[ ! -f /tmp/multilib_updates ]] ; then
    export windowtomoveto
    export wintomv=${pkgdir}
    wlm-fn_windowlabelandmove
  else
  rm /tmp/multilib_updates

  fi
  ls -R ${excludes} > FileList
  echo "There are system updates to check over and install"
else
  echo "Nothing new"
fi
echo
echo

## Script ends ##
